# csw/mgar/pkg/logstash/trunk/files/cswlogstash.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>

<service_bundle type='manifest' name='export'>
  
  <service name='network/csw/logstash' type='service' version='0'>
    
    <dependency name='fs' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/system/filesystem/local'/>
    </dependency>
    
    <dependency name='networking' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/milestone/network:default'/>
    </dependency>
    
    <exec_method name='start' type='method' exec='/opt/csw/bin/logstash agent -f %{logstash/config}' timeout_seconds='60'>
    <exec_method name='stop' type='method' exec=':kill' timeout_seconds='60'/>
    
    <instance name='default' enabled='false'>
        <property_group name='logstash' type='application'>
            <propval name='value_authorization' type='astring' value='solaris.smf.value.logstash'/>
            <propval name='config' type='astring' value='/etc/opt/csw/logstash/conf.d' />
        </property_group>
        <property_group name='startd' type='framework'>
          <propval name='ignore_error' type='astring' value='core,signal'/>
        </property_group>
        <property_group name='general' type='framework'>
          <propval name='action_authorization' type='astring' value='solaris.smf.manage.logstash'/>
        </property_group>
    </instance>

        <stability value='Evolving' />

        <template>
                <common_name>
                        <loctext xml:lang='C'>
                                Elasticsearch
                        </loctext>
                </common_name>
                <documentation>
                        <doc_link name='logstash'
                                uri='http://www.logstash.net' />
                </documentation>
        </template>

    
  </service>
  
</service_bundle>
