# csw/branches/user-yann/trunk/server/vsftpd/files/vsftpd.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<!--
William Pool (Puddle) 02/05
Service manifest for vsftpd
E-mail:
Modified for CSWmysql5 by Alex Moore 01/06
Modified for CSWvsftpd by Yann Rouillard 04/06
-->

<service_bundle type='manifest' name='vsftpd:vsftpd'>

<service
name='network/cswvsftpd'
type='service'
version='1'>
<create_default_instance enabled='false' />
<single_instance />

<dependency name='fs'
grouping='require_all'
restart_on='none'
type='service'>
<service_fmri value='svc:/system/filesystem/local' />
</dependency>

<dependency name='net'
grouping='require_all'
restart_on='none'
type='service'>
<service_fmri value='svc:/network/loopback' />
</dependency>

<exec_method
type='method'
name='start'
exec='/opt/csw/lib/svc/method/svc-vsftpd start'
timeout_seconds='-1'>
</exec_method>

<exec_method
type='method'
name='stop'
exec='/opt/csw/lib/svc/method/svc-vsftpd stop'
timeout_seconds='-1'>
</exec_method>

<exec_method
type='method'
name='restart'
exec='/opt/csw/lib/svc/method/svc-vsftpd restart'
timeout_seconds='-1'>
</exec_method>

</service>

</service_bundle>
