# csw/mgar/pkg/nagios/trunk/files/nagios.xml
<?xml version='1.0'?>
 <!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<!--

        Alex Moore 01/06
        Service manifest for nagios
 -->

 <service_bundle type='manifest' name='CSWnagios:nagios'>

 <service
         name='application/cswnagios'
         type='service'
         version='1'>

         <create_default_instance enabled='false' />

         <single_instance />

         <dependency name='filesystem'
                 grouping='require_all'
                 restart_on='none'
                 type='service'>
                 <service_fmri value='svc:/system/filesystem/local' />
         </dependency>

         <dependency
                 name='config-file'
                 grouping='require_any'
                 restart_on='refresh'
                 type='path'>
                 <service_fmri
                     value='file://localhost/opt/csw/nagios/etc/nagios.cfg' />
                </dependency>


         <dependency name='net'
                 grouping='require_all'
                 restart_on='none'
                 type='service'>
                 <service_fmri value='svc:/network/loopback' />
         </dependency>

         <dependency name='network'
                 grouping='require_all'
                 restart_on='none'
                 type='service'>
                 <service_fmri value='svc:/network/service' />
         </dependency>

         <dependency name='syslog'
                 grouping='optional_all'
                 restart_on='none'
                 type='service'>
                 <service_fmri value='svc:/system/system-log' />
         </dependency>

 	 <dependent
		 name='cswnagios_multi-user'
		 grouping='optional_all'
		 restart_on='none'>
		 <service_fmri
			value='svc:/milestone/multi-user' />
	 </dependent>

         <exec_method
                 type='method'
                 name='start'
                 exec='/opt/csw/lib/svc/method/svc-nagios start'
                 timeout_seconds='30'>
         </exec_method>

         <exec_method
                 type='method'
                 name='stop'
                 exec='/opt/csw/lib/svc/method/svc-nagios stop'
                 timeout_seconds='60'>
         </exec_method>

         <exec_method
                 type='method'
                 name='restart'
                 exec='/opt/csw/lib/svc/method/svc-nagios restart'
                 timeout_seconds='60'>
         </exec_method>

 </service>

 </service_bundle>
