# csw/mgar/server/jetty6/files/cswjetty6.xml
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>

<service_bundle type='manifest' name='CSWjetty6:jetty6'>
    <service name='network/cswjetty6' type='service' version='1'>

        <create_default_instance enabled='false'/>

        <dependency
            name='multi-user'
            grouping='require_all'
            restart_on='error'
            type='service'>
            <service_fmri value='svc:/milestone/multi-user'/>
        </dependency>
        
        <dependency
            name='nfs'
            grouping='optional_all'
            restart_on='error'
            type='service'>
            <service_fmri value='svc:/network/nfs/client'/>
        </dependency>

        <method_context working_directory="/">
            <method_credential user="jetty" group="jetty"/>

            <method_environment>
                <envvar name="JAVA_HOME" value="/usr/jdk/latest"/>
            </method_environment>
        </method_context>
        
        <exec_method
            name='start'
            type='method'
            exec='/opt/csw/lib/svc/method/svc-cswjetty6 start'
            timeout_seconds='30'>
        </exec_method>
        
        <exec_method
            name='stop'
            type='method'
            exec='/opt/csw/lib/svc/method/svc-cswjetty6 stop'
            timeout_seconds='30'>
        </exec_method>

        <exec_method
            name='refresh'
            type='method'
            exec='/opt/csw/lib/svc/method/svc-cswjetty6 refresh'
            timeout_seconds='30'>
        </exec_method>        

        <stability value='Unstable'/>
        
        <template>
            <common_name>
                <loctext xml:lang='C'>Jetty6 Servlet Container</loctext>
            </common_name>
<!--
            <documentation>
                <manpage 
                    title='postfix' 
                    section='1'
                    manpath='/opt/csw/share/man'/>
            </documentation>
-->
        </template>
    </service>
</service_bundle>
