# csw/branches/modularpkg/pkg/net/cups/files/cups.xml
<?xml version='1.0'?>
 <!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<!--

        William Pool (Puddle) 02/05
        Service manifest for CUPS
        E-mail: puddle@flipmotion.com
	Modified for CSWcups by Markus Mayer <mmayer@blastwave.org>, 2006-01-04
 -->

 <service_bundle type='manifest' name='CSWcups:cups'>

 <service
         name='application/print/cswcups'
         type='service'
         version='1'>
         <create_default_instance enabled='false' />
         <single_instance />


         <dependency name='fs'
                 grouping='require_all'
                 restart_on='none'
                 type='service'>
                 <service_fmri value='svc:/system/filesystem/local' />
         </dependency>

         <dependency name='net'
                 grouping='require_all'
                 restart_on='none'
                 type='service'>
                 <service_fmri value='svc:/network/loopback' />
         </dependency>

         <exec_method
                 type='method'
                 name='start'
                 exec='/opt/csw/lib/svc/method/svc-cups start'
                 timeout_seconds='-1'>
         </exec_method>

         <exec_method
                 type='method'
                 name='stop'
                 exec='/opt/csw/lib/svc/method/svc-cups stop'
                 timeout_seconds='-1'>
         </exec_method>

         <exec_method
                 type='method'
                 name='reload'
                 exec='/opt/csw/lib/svc/method/svc-cups reload'
                 timeout_seconds='-1'>
         </exec_method>

         <exec_method
                 type='method'
                 name='restart'
                 exec='/opt/csw/lib/svc/method/svc-cups restart'
                 timeout_seconds='-1'>
         </exec_method>

         <exec_method
                 type='method'
                 name='status'
                 exec='/opt/csw/lib/svc/method/svc-cups status'
                 timeout_seconds='-1'>
         </exec_method>

 </service>

 </service_bundle>
