# csw/mgar/pkg/fusioninventory-agent/trunk/files/fusioninventory-agent.cron
40 3 * * * . /etc/opt/csw/default/fusioninventory-agent && /opt/csw/bin/fusioninventory-agent $OPTIONS > /dev/null 2>&1
