# csw/mgar/pkg/lang-python/carbon/trunk/files/cswcarbon-relay.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>

<service_bundle type='manifest' name='export'>
  
  <service name='network/csw/carbon-relay' type='service' version='0'>
    
    <dependency name='fs' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/system/filesystem/local'/>
    </dependency>
    
    <dependency name='networking' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/milestone/network:default'/>
    </dependency>
    
    <exec_method name='start' type='method' exec='/opt/csw/bin/carbon-relay --config=%{carbonrelay/config} start' timeout_seconds='60'/>
    <exec_method name='stop' type='method' exec='/opt/csw/bin/carbon-relay --config=%{carbonrelay/config} stop' timeout_seconds='60'/>
    
    <instance name='default' enabled='false'>
        <property_group name='carbonrelay' type='application'>
            <propval name='value_authorization' type='astring' value='solaris.smf.value.graphite'/>
            <propval name='config' type='astring' value='/etc/opt/csw/graphite/relay-rules.conf' />
        </property_group>
        <property_group name='startd' type='framework'>
          <propval name='ignore_error' type='astring' value='core,signal'/>
          <propval name='duration' type='astring' value='contract' />
        </property_group>
        <property_group name='general' type='framework'>
          <propval name='action_authorization' type='astring' value='solaris.smf.manage.graphite'/>
        </property_group>
    </instance>

        <stability value='Evolving' />

        <template>
                <common_name>
                        <loctext xml:lang='C'>
                                Graphite
                        </loctext>
                </common_name>
                <documentation>
                        <doc_link name='graphite'
                                uri='http://graphite.wikidot.com/carbon' />
                </documentation>
        </template>

    
  </service>
  
</service_bundle>
