# csw/mgar/pkg/exim/trunk/files/cswexim.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<!--
Service manifest for cswexim
Autogenerated by CSWinitsmf
2008-04-21 Peter Bonivart

Modifications since autogenration:
1. Depend on sendmail-client
-->

<service_bundle type='manifest' name='cswexim:cswexim'>

<service
  name='network/cswexim'
  type='service'
  version='1'>
  <create_default_instance enabled='false' />
  <single_instance />

  <dependency name='fs'
    grouping='require_all'
    restart_on='none'
    type='service'>
    <service_fmri value='svc:/system/filesystem/local' />
  </dependency>

  <dependency name='net'
    grouping='require_all'
    restart_on='none'
    type='service'>
    <service_fmri value='svc:/network/loopback' />
  </dependency>

  <dependency name='sendmail'
    grouping='require_all'
    restart_on='none'
    type='service'>
    <service_fmri value='svc:/network/sendmail-client' />
  </dependency>

  <exec_method
    type='method'
    name='start'
    exec='/var/opt/csw/svc/method/svc-cswexim start'
    timeout_seconds='120'>
  </exec_method>

  <exec_method
    type='method'
    name='stop'
    exec='/var/opt/csw/svc/method/svc-cswexim stop'
    timeout_seconds='60'>
  </exec_method>

  <exec_method
    type='method'
    name='restart'
    exec='/var/opt/csw/svc/method/svc-cswexim restart'
    timeout_seconds='180'>
  </exec_method>

</service>

</service_bundle>
