# csw/mgar/pkg/nsca/trunk/files/nsca.xml
<?xml version='1.0'?>
 <!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<!--

        Juergen Arndt (opencsw.org) 02/09
        Service manifest for nsca
 -->

 <service_bundle type='manifest' name='nsca:nrpe'>

 <service
         name='network/cswnsca'
         type='service'
         version='1'>
         <create_default_instance enabled='false' />
         <single_instance />


         <dependency name='fs'
                 grouping='require_all'
                 restart_on='none'
                 type='service'>
                 <service_fmri value='svc:/system/filesystem/local' />
         </dependency>

         <dependency
                 name='config-file'
                 grouping='require_any'
                 restart_on='refresh'
                 type='path'>
                 <service_fmri
                     value='file://localhost/opt/csw/nagios/etc/nsca.cfg' />
         </dependency>

	<dependency
		name='nagios-command-file'
		grouping='require_any'
		restart_on='none'
		type='path'>
		<service_fmri value='file://localhost/opt/csw/nagios/var/rw/nagios.cmd' />
	</dependency>


         <dependency name='net'
                 grouping='require_all'
                 restart_on='none'
                 type='service'>
                 <service_fmri value='svc:/network/loopback' />
         </dependency>

         <dependency name='network'
                 grouping='require_all'
                 restart_on='none'
                 type='service'>
                 <service_fmri value='svc:/network/service' />
         </dependency>

         <dependency name='syslog'
                 grouping='optional_all'
                 restart_on='none'
                 type='service'>
                 <service_fmri value='svc:/system/system-log' />
         </dependency>

         <exec_method
                 type='method'
                 name='start'
                 exec='/opt/csw/lib/svc/method/svc-nsca start'
                 timeout_seconds='-1'>
         </exec_method>

         <exec_method
                 type='method'
                 name='stop'
                 exec='/opt/csw/lib/svc/method/svc-nsca stop'
                 timeout_seconds='-1'>
         </exec_method>


 </service>

 </service_bundle>
